# buildroot/board/gsd/boson-s3-dev/skeleton/etc/systemd/system/led_off.service
[Unit]
Description=LED Toggle service

[Service]
ExecStart=/usr/sbin/led_off.sh

[Install]
WantedBy=multi-user.target